# ai-chatbot.service
# Systemd user service for AI Chatbot
# Install to: ~/.config/systemd/user/ai-chatbot.service
# Enable with: systemctl --user enable ai-chatbot
# Start with: systemctl --user start ai-chatbot
# View logs: journalctl --user -u ai-chatbot -f

[Unit]
Description=AI Chatbot - Personal Gemini-powered assistant
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/src/ai-chatbot

# Build frontend before starting (npm install + build)
ExecStartPre=/bin/bash -c 'cd %h/src/ai-chatbot/web && npm install && npm run build'

ExecStart=%h/src/ai-chatbot/.venv/bin/gunicorn --bind 0.0.0.0:${PORT} --workers ${GUNICORN_WORKERS:-2} --timeout ${GUNICORN_TIMEOUT:-300} "src.app:create_app()"
Restart=on-failure
RestartSec=5
Environment=PATH=%h/src/ai-chatbot/.venv/bin:/usr/local/bin:/usr/bin:/bin

# Load environment from .env file
EnvironmentFile=%h/src/ai-chatbot/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
