# systemd/ai-chatbot-file-cleanup.timer
# Systemd user timer for AI Chatbot file retention cleanup
# Install to: ~/.config/systemd/user/ai-chatbot-file-cleanup.timer
# Enable with: systemctl --user enable --now ai-chatbot-file-cleanup.timer

[Unit]
Description=Daily cleanup of expired attachments (videos 7d, images/files 30d)

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
