# systemd/ai-chatbot-checkpoint-cleanup.service
# Systemd user service for AI Chatbot checkpoint cleanup
# Install to: ~/.config/systemd/user/ai-chatbot-checkpoint-cleanup.service
# This service is triggered by ai-chatbot-checkpoint-cleanup.timer

[Unit]
Description=AI Chatbot - Checkpoint cleanup (delete expired LangGraph state)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/src/ai-chatbot
ExecStart=%h/src/ai-chatbot/.venv/bin/python scripts/cleanup_checkpoints.py
Environment=PATH=%h/src/ai-chatbot/.venv/bin:/usr/local/bin:/usr/bin:/bin

# Load environment from .env file (for database paths, TTL config)
EnvironmentFile=%h/src/ai-chatbot/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Don't restart on failure (timer will retry next cycle)
Restart=no

[Install]
WantedBy=default.target
